# ops/roles/authark/templates/authark.service.j2
[Unit]
Description=Authark Service

[Service]
User=authark
Group=authark
WorkingDirectory=/opt/authark/authark
ExecStart=/opt/authark/env/bin/python -m authark
Restart=on-failure